# configs/osx.crontab
@daily sudo bash /home/josh/projects/personal-computer-provisioning/cron/daily.sh
@weekly sudo bash /home/josh/projects/personal-computer-provisioning/cron/weekly.sh
@hourly sudo bash /home/josh/projects/personal-computer-provisioning/cron/hourly.sh
