# configs/osx.crontab
@daily /Users/joshrogan/projects/personal-computer-provisioning/cron/daily.sh
@hourly /Users/joshrogan/projects/personal-computer-provisioning/cron/hourly.sh
@weekly /Users/joshrogan/projects/personal-computer-provisioning/cron/weekly.sh
@reboot /Users/joshrogan/projects/personal-computer-provisioning/cron/reboot.sh
* * * * * /Users/joshrogan/projects/personal-computer-provisioning/cron/hourly.sh
* * * * * /Users/joshrogan/projects/personal-computer-provisioning/cron/cron.all.log
